# vlekt.service
[Unit]
Description=VLEKT PRO - Gestionale Nutrizionale
After=network.target

[Service]
Type=simple
# Cambia con il tuo utente (es. root, ubuntu, www-data)
User=root
Group=root
WorkingDirectory=/opt/vlekt
Environment="VLEKT_LICENSE_KEY=tua_chiave_qui"
ExecStart=/opt/vlekt/venv/bin/python3 -m streamlit run app2.py --server.port 8501 --server.address 0.0.0.0 --server.headless true --browser.gatherUsageStats false
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
